# wsl.timer
[Unit]
Description=Runs wsl generator every Mon

[Timer]
OnCalendar=Mon *-*-* 01:00:00
Unit=wsl.service

[Install]
WantedBy=multi-user.target
